# com.jordan.yad2watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jordan.yad2watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jordanperez/Desktop/Projects/yad2-watcher/watcher.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jordanperez/Desktop/Projects/yad2-watcher</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jordanperez/Desktop/Projects/yad2-watcher/data/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jordanperez/Desktop/Projects/yad2-watcher/data/launchd.log</string>
</dict>
</plist>
